# rootfs/etc/services.d/mdns/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start mDNS repeater service
# ==============================================================================
interface=$(awk '{for (i=1; i<=NF; i++) if ($i~/dev/) print $(i+1)}' <<< "$(ip route get 8.8.8.8)")

exec mdns-repeater -f ${interface} hassio