# rootfs/etc/services.d/mdns/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start mDNS repeater service
# ==============================================================================

exec mdns-repeater -f hassio
