# rootfs/etc/services.d/pimd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start pimd
# ==============================================================================

exec pimd -f
